# .config/systemd/user/maven-prune.timer
[Unit]
Description=Timer to prune the maven storage

[Timer]
OnStartupSec=1 minute
OnUnitInactiveSec=6 hours
