# bin/dump-ext.cron
#!/bin/bash
#

# Get the directory that contains the real command file.
cron_script_file="$0"
script_dir=$(dirname $(readlink -e "$cron_script_file"))

awk '
$3 ~ /ext./ && $5 == 1 {
	print $1
}
' /etc/fstab | xargs -I{} -P 0 "$script_dir/dump-ext" {} "/var/backups/$(uname -n)"
